# sentinel-alpha@.service
[Unit]
Description=Sentinel Alpha WEEX Trading Bot (%I)
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/sentinel-alpha
Environment="PATH=/root/sentinel-alpha/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/root/sentinel-alpha/venv/bin/python3 /root/sentinel-alpha/live_trading_bot.py --skip-prompt --symbol %i
Restart=always
RestartSec=10
# Consolidate ALL instance logs to one file for easier monitoring
StandardOutput=append:/root/sentinel-alpha/logs/combined_bots.log
StandardError=append:/root/sentinel-alpha/logs/combined_bots.log

# Resource limits
MemoryLimit=1G
CPUQuota=40%

# Restart policy
StartLimitBurst=5
StartLimitIntervalSec=300

[Install]
WantedBy=multi-user.target
